# modules/daemon/templates/systemd.service.erb
[Unit]
Description=<%= @name %>
Before=<%= @daemon_before.join(' ') unless @daemon_before.nil? %>
After=network.target <%= @daemon_after.join(' ') unless @daemon_after.nil? %>

[Service]
Type=simple
<% if @pre_command %>ExecStartPre=<%= @pre_command %><% end %>
ExecStart=<%= @binary %> <%= @args %>
<% if @post_command %>ExecStartPost=<%= @post_command %><% end %>
User=<%= @user %>
PermissionsStartOnly=true
<% if @nice %>Nice=<%= @nice %><% end %>
<% if @oom_score_adjust %>OOMScoreAdjust=<%= @oom_score_adjust %><% end %>
TimeoutStopSec=<%= @stop_timeout %>
KillMode=mixed
Restart=always
<% if @limit_nofile %>LimitNOFILE=<%= @limit_nofile %><% end %>
<% if @core_dump %>LimitCORE=infinity<% end %>
<% @env.each do |name, value| -%>
Environment=<%= name %>=<%= value %>
<% end %>

[Install]
WantedBy=multi-user.target
