# modules/daemon/templates/systemd.service.erb
[Unit]
Description=<%= @name %>
After=network.target

[Service]
Type=simple
ExecStart=<%= @binary %> <%= @args %>
User=<%= @user %>
<% if @nice %>Nice=<%= @nice %><% end %>
TimeoutStopSec=<%= @stop_timeout %>
KillMode=mixed
Restart=always

[Install]
WantedBy=multi-user.target
