# modules/daemon/templates/sysvinit.sh.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:           <%= @name %>
# Required-Start:     $local_fs $syslog
# Required-Stop:      $local_fs $syslog
# Default-Start:      2 3 4 5
# Default-Stop:       0 1 6
# Short-Description:  <%= @name %>
### END INIT INFO

NAME=<%= @name %>
DESC=<%= @name %>
DAEMON='<%= @binary %>'
DAEMON_ARGS='<%= @args %>'
PIDFILE='/var/run/<%= @name %>.pid'
USER='<%= @user %>'
<% unless @env.empty? %>
<% @env.each do |name, value| %>
export <%= name %>\='<%= value %>'
<% end %>
<% end %>
<% @final_signal = @sysvinit_kill ? 'KILL' : 'TERM' -%>

set -e
. /lib/lsb/init-functions

case "${1}" in
	start)
		log_daemon_msg "Starting ${DESC}" "${NAME}"
		<% if @core_dump %>
		ulimit -c unlimited
		<%- end %>
		<% if @limit_nofile %>
		ulimit -n <%= @limit_nofile == 'unlimited' ? 'unlimited' : @limit_nofile %>
		<%- end %>
		<% if @limit_fsize %>
		ulimit -f <%= @limit_fsize == 'unlimited' ? 'unlimited' : @limit_fsize %>
		<%- end %>
		<% if @limit_cpu %>
		ulimit -t  <%= @limit_cpu == 'unlimited' ? 'unlimited' : @limit_cpu %>
		<%- end %>
		<% if @limit_as %>
		ulimit -v  <%= @limit_as == 'unlimited' ? 'unlimited' : @limit_as %>
		<%- end %>
		<% if @limit_rss %>
		ulimit -m  <%= @limit_rss == 'unlimited' ? 'unlimited' : @limit_rss %>
		<%- end %>
		<% if @limit_nproc %>
		ulimit -u  <%= @limit_nproc == 'unlimited' ? 'unlimited' : @limit_nproc %>
		<%- end %>
		<% if @pre_command %>
		<%= @pre_command %>
		<%- end %>
		<% if @runtime_directory %>
		mkdir -p /run/<%= @runtime_directory %>
		<% if @runtime_directory_mode %>
		chmod <%= @runtime_directory_mode %> /run/<%= @runtime_directory %>
		<%- end %>
		<%- end %>
		if (start-stop-daemon --start --oknodo --make-pidfile --background --pidfile "${PIDFILE}" --chuid "${USER}" --startas "${DAEMON}" -- ${DAEMON_ARGS}); then
		  timeout --signal=9 1 bash -c "while ! (test -e '${PIDFILE}'); do sleep 0.05; done";
		  PID=$(cat "${PIDFILE}")
			<% if @nice %>
				renice -n <%= @nice %> -p ${PID}
			<%- end %>
			<% if @oom_score_adjust %>
				echo <%= @oom_score_adjust %> > /proc/${PID}/oom_score_adj
			<%- end %>
			log_end_msg 0
			<% if @post_command %>
			<%= @post_command %>
			<%- end %>
		else
			log_end_msg 1
		fi
	;;
	stop)
		log_daemon_msg "Stopping ${DESC}" "${NAME}"
		if (start-stop-daemon --stop --oknodo --retry TERM/<%= @stop_timeout %>/<%= @final_signal %>/<%= @stop_timeout %> --pidfile "${PIDFILE}" --user "${USER}"); then
			rm -f "${PIDFILE}"
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	status)
		status_of_proc -p "${PIDFILE}" "${DAEMON}" "${NAME}"
	;;
	restart)
		${0} stop
		${0} start
	;;
	*)
		echo "Usage: ${0} {start|stop|status|restart}" >&2
	;;
esac

exit 0
