# modules/daemon/templates/spec/foo-before.service
[Unit]
Description=Foo Before Service

[Service]
Type=simple
ExecStart=/usr/local/bin/my-program-before
