# modules/daemon/templates/spec/foo-after.service
[Unit]
Description=Foo After Service

[Service]
Type=simple
ExecStart=/usr/local/bin/my-program-after
